# initrd/musb-enable.service
[Unit]
Description=Enable MUSB gadgets

[Service]
Type=simple
ExecStart=/var/lib/misc/enable-musb.sh

[Install]
WantedBy=multi-user.target
